# com.neo.trading-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.neo.trading-dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/mac/.openclaw/workspace-neo/dashboard/app.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/mac/.openclaw/workspace-neo/dashboard</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/mac/.openclaw/workspace-neo/logs/dashboard.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mac/.openclaw/workspace-neo/logs/dashboard_error.log</string>
</dict>
</plist>
